# ops/systemd/drama-mud.service
[Unit]
Description=Drama MUD server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__APP_USER__
Group=__APP_USER__
WorkingDirectory=__APP_ROOT__/current
Environment=NODE_ENV=production
Environment=PATH=/usr/local/bin:/usr/bin:/bin
EnvironmentFile=-__APP_ROOT__/shared/drama-mud.env
ExecStart=/usr/bin/env pnpm --dir __APP_ROOT__/current --filter @drama-mud/server start
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
